# favor.service
[Unit]
Description=Favor AI Companion
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/node /root/whatsapp-bot/favor.js
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/root
Environment=ANTHROPIC_API_KEY=your-anthropic-api-key-here
Environment=BRAVE_API_KEY=your-brave-api-key-here
Environment=TMPDIR=/tmp
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
WorkingDirectory=/root/whatsapp-bot

[Install]
WantedBy=default.target
